# deploy/supervisor.conf
# Supervisor 配置 - 后端服务
[program:zulin_backend]
command=/var/www/zulin/venv/bin/gunicorn -c gunicorn.conf.py app.main:app
directory=/var/www/zulin/backend
autostart=true
autorestart=true
stderr_logfile=/var/log/zulin/backend_err.log
stdout_logfile=/var/log/zulin/backend_out.log
user=www-data
group=www-data
environment=PATH="/var/www/zulin/venv/bin",PYTHONPATH="/var/www/zulin/backend"

# Supervisor 配置 - 前端构建（可选，如果需要自动构建）
# [program:zulin_frontend_build]
# command=npm run build
# directory=/var/www/zulin/frontend
# autostart=false
# autorestart=false
